# Daily/scheduler/com.india-ts.brooks_reversal_simple.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.india-ts.brooks_reversal_simple</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.12/bin/python3</string>
        <string>/Users/maverick/PycharmProjects/India-TS/Daily/utils/brooks_reversal_scheduler.py</string>
    </array>
    
    <!-- Run every 30 minutes (1800 seconds) -->
    <key>StartInterval</key>
    <integer>1800</integer>
    
    <!-- Log file paths -->
    <key>StandardErrorPath</key>
    <string>/Users/maverick/PycharmProjects/India-TS/Daily/logs/brooks_reversal_error.log</string>
    <key>StandardOutPath</key>
    <string>/Users/maverick/PycharmProjects/India-TS/Daily/logs/brooks_reversal.log</string>
    
    <!-- Run when loaded -->
    <key>RunAtLoad</key>
    <true/>
    
    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
    
    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/maverick/PycharmProjects/India-TS/Daily</string>
    
    <!-- Don't keep alive - let StartInterval handle the scheduling -->
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>