# Daily/scheduler/com.india-ts.update_cnc_stoploss.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.india-ts.update_cnc_stoploss</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/maverick/PycharmProjects/India-TS/.venv/bin/python</string>
        <string>/Users/maverick/PycharmProjects/India-TS/Daily/scripts/update_cnc_stoploss.py</string>
        <!-- Force placement outside market hours -->
        <string>--force-place</string>
        <!-- Force refresh data from API -->
        <string>--refresh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>14</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardErrorPath</key>
    <string>/Users/maverick/PycharmProjects/India-TS/logs/update_cnc_stoploss_error.log</string>
    <key>StandardOutPath</key>
    <string>/Users/maverick/PycharmProjects/India-TS/logs/update_cnc_stoploss.log</string>
    <key>RunAtLoad</key>
    <false/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/maverick/PycharmProjects/India-TS</string>
</dict>
</plist>