# Daily/scheduler/plists/com.india-ts.regime_data_updater_10min.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.india-ts.regime_data_updater_10min</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/maverick/PycharmProjects/India-TS/Daily/Market_Regime/update_regime_data.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/maverick/PycharmProjects/India-TS/Daily/Market_Regime</string>
    
    <key>StartCalendarInterval</key>
    <array>
        <!-- Run at :10, :20, :30, :40, :50 past each hour from 9 AM to 3 PM -->
        <!-- 9:10, 9:20, 9:30, 9:40, 9:50 -->
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>10</integer></dict>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>20</integer></dict>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>40</integer></dict>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>50</integer></dict>
        
        <!-- 10:00, 10:10, 10:20, 10:30, 10:40, 10:50 -->
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>10</integer></dict>
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>20</integer></dict>
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>40</integer></dict>
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>50</integer></dict>
        
        <!-- 11:00 through 11:50 -->
        <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>10</integer></dict>
        <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>20</integer></dict>
        <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>40</integer></dict>
        <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>50</integer></dict>
        
        <!-- 12:00 through 12:50 -->
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>10</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>20</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>40</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>50</integer></dict>
        
        <!-- 1:00 PM through 1:50 PM -->
        <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>10</integer></dict>
        <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>20</integer></dict>
        <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>40</integer></dict>
        <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>50</integer></dict>
        
        <!-- 2:00 PM through 2:50 PM -->
        <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>10</integer></dict>
        <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>20</integer></dict>
        <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>40</integer></dict>
        <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>50</integer></dict>
        
        <!-- 3:00 PM through 3:30 PM -->
        <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>10</integer></dict>
        <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>20</integer></dict>
        <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>30</integer></dict>
    </array>
    
    <key>StandardOutPath</key>
    <string>/Users/maverick/PycharmProjects/India-TS/Daily/logs/regime_data_updater.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/maverick/PycharmProjects/India-TS/Daily/logs/regime_data_updater_error.log</string>
    
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>