# Daily/scheduler/plists/com.india-ts.hourly-breakout-alerts.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.india-ts.hourly-breakout-alerts</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/maverick/PycharmProjects/India-TS/.venv/bin/python</string>
        <string>/Users/maverick/PycharmProjects/India-TS/Daily/alerts/hourly_breakout_alert_service.py</string>
        <string>-u</string>
        <string>Sai</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/maverick/PycharmProjects/India-TS</string>
    
    <key>StandardOutPath</key>
    <string>/Users/maverick/PycharmProjects/India-TS/Daily/logs/alerts_hourlybo/hourly_breakout_alerts.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/maverick/PycharmProjects/India-TS/Daily/logs/alerts_hourlybo/hourly_breakout_alerts_error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/maverick/PycharmProjects/India-TS</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    
    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>