# Daily/scheduler/plists/com.india-ts.weekly_backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.india-ts.weekly_backup</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/maverick/PycharmProjects/India-TS/Daily/bin/weekly_backup.sh</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <!-- Run every Sunday at 3:00 AM -->
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>WorkingDirectory</key>
    <string>/Users/maverick/PycharmProjects/India-TS</string>
    
    <key>StandardOutPath</key>
    <string>/Users/maverick/PycharmProjects/India-TS/Daily/logs/weekly_backup.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/maverick/PycharmProjects/India-TS/Daily/logs/weekly_backup_error.log</string>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <false/>
    
    <key>ThrottleInterval</key>
    <integer>3600</integer>
</dict>
</plist>