# Daily/scheduler/plists/com.india-ts.fix_plists_on_startup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.india-ts.fix_plists_on_startup</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/maverick/PycharmProjects/India-TS/Daily/scripts/fix_brooks_plist.sh</string>
    </array>
    
    <!-- Run at system startup -->
    <key>RunAtLoad</key>
    <true/>
    
    <!-- Log file paths -->
    <key>StandardErrorPath</key>
    <string>/Users/maverick/PycharmProjects/India-TS/Daily/logs/plist_fix_error.log</string>
    <key>StandardOutPath</key>
    <string>/Users/maverick/PycharmProjects/India-TS/Daily/logs/plist_fix_output.log</string>
    
    <!-- Don't keep running -->
    <key>KeepAlive</key>
    <false/>
    
    <!-- Delay startup by 30 seconds to ensure system is ready -->
    <key>StartInterval</key>
    <integer>86400</integer>
    
    <!-- Exit after running once -->
    <key>LaunchOnlyOnce</key>
    <true/>
</dict>
</plist>