# Daily/scheduler/plists/com.india-ts.hourly-short-tracker-service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.india-ts.hourly-short-tracker-service</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/python3</string>
		<string>/Users/maverick/PycharmProjects/India-TS/Daily/services/hourly_short_tracker_service.py</string>
		<string>--user</string>
		<string>Sai</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>StandardErrorPath</key>
	<string>/Users/maverick/PycharmProjects/India-TS/Daily/logs/hourly_short_tracker_service_error.log</string>
	<key>StandardOutPath</key>
	<string>/Users/maverick/PycharmProjects/India-TS/Daily/logs/hourly_short_tracker_service.log</string>
	<key>WorkingDirectory</key>
	<string>/Users/maverick/PycharmProjects/India-TS/Daily</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>PYTHONPATH</key>
		<string>/Users/maverick/PycharmProjects/India-TS</string>
	</dict>
</dict>
</plist>